# card_skeleton/etc/update-motd.d/10-uname
==================================================
SPIS - Simple PI SendSpin, version ${SPIS_VERSION}
==================================================
